# Engram/org.klevan.Engram.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Privileged helper daemon for one-click CLI install (ADR 0022).
         Registered by the app via SMAppService.daemon(plistName:); must live at
         Contents/Library/LaunchDaemons/ inside the app bundle. The program is the
         bundled engram CLI run with the hidden _helper-daemon subcommand, which
         vends the org.klevan.Engram.helper Mach service. -->
    <key>Label</key>
    <string>org.klevan.Engram.helper</string>
    <key>BundleProgram</key>
    <string>Contents/Helpers/engram</string>
    <key>ProgramArguments</key>
    <array>
        <string>Contents/Helpers/engram</string>
        <string>_helper-daemon</string>
    </array>
    <key>MachServices</key>
    <dict>
        <key>org.klevan.Engram.helper</key>
        <true/>
    </dict>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>org.klevan.Engram</string>
    </array>
</dict>
</plist>
